# lesson-34/provisioning/roles/nfs-client/files/srv-nfs.automount
[Unit]
Description=OTUS mount NFS3
Requires=network-online.target

[Automount]
Where=/srv/nfs
TimeoutIdleSec=60

[Install]
WantedBy = multi-user.target
