# lesson-08/p1/provisioning/roles/otus-p1/files/mon-log.timer
[Unit]
Description=Monitoring log timer

[Timer]
OnCalendar=*:*:0/30
AccuracySec=1s

[Install]
WantedBy=multi-user.target
